# client/client.service
[Unit]
Description=NYILSRV Client
After=network.target

[Service]
Type=simple
ExecStart=/opt/nyilsrv-client/venv/bin/python /opt/nyilsrv-client/client.py
Restart=always
WorkingDirectory=/opt/nyilsrv-client
User=root
Environment=PATH=/opt/nyilsrv-client/venv/bin:/usr/bin:/bin

[Install]
WantedBy=multi-user.target
